# roles/fedora_container_sync/templates/fedora-container-sync-external.service.jinja2
# /etc/systemd/system/fedora-container-sync-external.service
[Unit]
Description=fedora-container-sync-external start service
After=network-online.target
[Service]
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/fedora-container-sync-external.sh
Type=oneshot
TimeoutSec=0
