# roles/fedora_cloud_sync/templates/fedora-cloud-sync-external.timer.jinja2
[Unit]
Description=fedora-cloud-sync-external timer
[Timer]
# fedora-cloud-sync-external will start 23:00, a hour later
OnCalendar=*-*-* 22:00:00
# Persistent=True triggers the service immediately if it missed the last start time
Persistent=True
Unit=fedora-cloud-sync-external.service
[Install]
WantedBy=timers.target
