# roles/fcos_sync/templates/fcos-sync-external.timer.jinja2
[Unit]
Description=fcos-sync-external timer
[Timer]
# fcos-sync-external will start 23:00, a hour later
OnCalendar=*-*-* 22:00:00
# Persistent=True triggers the service immediately if it missed the last start time
Persistent=True
Unit=fcos-sync-external.service
[Install]
WantedBy=timers.target
